# debian/postinst
#!/bin/sh

set -e

case "$1" in
	configure)
		getent group uffd >/dev/null 2>&1 || addgroup --system uffd
		adduser --system --home /var/lib/uffd --quiet uffd --ingroup uffd || true

		chown -R uffd:uffd /var/lib/uffd
		chmod 0770 /var/lib/uffd

		python3 <<EOF
import secrets
cfg = open('/etc/uffd/uffd.cfg', 'r').read()
cfg = cfg.replace('\n#SECRET=autogenerated by postinst script\n',
                  '\nSECRET="'+secrets.token_hex(128)+'"\n', 1)
open('/etc/uffd/uffd.cfg', 'w').write(cfg)
EOF
		chown root:uffd /etc/uffd/uffd.cfg
		chmod 0640 /etc/uffd/uffd.cfg

		invoke-rc.d uwsgi restart uffd
	;;

	abort-upgrade|abort-remove|abort-deconfigure)
	;;

	*)
		echo "postinst called with unknown argument \`$1'" >&2
		exit 1
	;;
esac

#DEBHELPER#

exit 0
